# deb/install/prerm
#!/bin/bash

systemctl stop waggle-internet-share.service || true
systemctl disable waggle-internet-share.service || true
