# deb/install/prerm
#!/bin/bash

systemctl stop waggle-internet-share.service || true
systemctl disable waggle-interface-share.service || true
